# server/certbot-renew-alert.service
[Unit]
Description=Alert on certbot renewal failure
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/certbot_renew_alert.sh
